# root_folder/etc/cron.d/update_pckg
#!/bin/bash
apt update
apt upgrade > /var/log/update_script.log
